# roles/ap/files/rc.local
iptables-restore < /etc/iptables.ipv4
exit 0